# deploy/systemd/newsromania-app.service
# newsromania-app.service — stack-ul de producție (docker compose --profile app).
# Unitate systemd de UTILIZATOR (systemctl --user), peste Docker ROOTLESS
# (docker.service din managerul de utilizator).
#
# oneshot + RemainAfterExit: unitatea doar pornește/oprește compose;
# containerele au propriul `restart: unless-stopped`, deci supravegherea
# proceselor rămâne la dockerd, nu la systemd.
#
# `systemctl --user enable newsromania-app` + linger (deja activ pentru
# newsagent) ⇒ stack-ul pornește automat la bootul serverului, fără sesiune
# deschisă. NU se activează la instalare — decizia aparține pasului de
# integrare / ownerului (vezi deploy/DEPLOY.md).
# Instalare: deploy/systemd/install-user-units.sh

[Unit]
Description=NewsRomania — aplicație producție (compose profile app, :3100)
Wants=network-online.target
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/newsagent/workspace/newsromania
Environment=DOCKER_HOST=unix:///run/user/1004/docker.sock
ExecStart=/bin/bash -lc 'export PATH="$$HOME/bin:$$PATH"; cd /home/newsagent/workspace/newsromania && docker compose --profile app up -d'
ExecStop=/bin/bash -lc 'export PATH="$$HOME/bin:$$PATH"; cd /home/newsagent/workspace/newsromania && docker compose --profile app down'
# Prima pornire poate aștepta healthcheck-urile postgres/redis.
TimeoutStartSec=300
TimeoutStopSec=120

[Install]
WantedBy=default.target
