# deploy/systemd/newsromania-ingest.service
# newsromania-ingest.service — o rulare a workerului de ingestie RSS.
# Unitate systemd de UTILIZATOR (systemctl --user). Fără secrete aici:
# `npx payload run` încarcă .env din rădăcina proiectului.
# Instalare: deploy/systemd/install-user-units.sh

[Unit]
Description=NewsRomania — ingestie RSS (agregare știri)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/newsagent/workspace/newsromania
# bash -lc ca să preia PATH-ul utilizatorului (node via ~/bin, regula CLAUDE.md).
ExecStart=/bin/bash -lc 'export PATH="$HOME/bin:$PATH"; cd /home/newsagent/workspace/newsromania && npx payload run scripts/worker/ingest.mjs'
# Timerul pornește la fiecare 5 min și procesează un singur lot rotativ (~10
# feeduri), deci o rulare se termină în zeci de secunde. Workerul se
# auto-limitează la 5 min (RUN_CAP_MS); systemd taie la 290s ca plasă de
# siguranță, SUB intervalul de 5 min, ca o rulare blocată să nu se suprapună
# cu următoarea.
TimeoutStartSec=290
# VPS partajat: prioritate redusă CPU/IO.
Nice=10
IOSchedulingClass=idle
